# private_dot_config/systemd/user/signal-i3status-rs.service
[Service]  
Type=oneshot
ExecStart=/usr/bin/pkill --signal "SIGRTMIN+1" "i3status-rs"

# vi: ft=systemd
